# crontab
*/5 * * * * cd /app && npm run start >> /var/log/cron.log 2>&1
